# ta_coe/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

if [[ $(bashio::config "debug") == "true" ]]; then
    export DEBUG=1
    bashio::log.info "Debug mode enabled"
fi

exec python3 /opt/ta-coe/main.py
